# cron/hta-cloud.crontab
# Example of job definition:
# .---------------- minute (0 - 59)
# |  .------------- hour (0 - 23)
# |  |  .---------- day of month (1 - 31)
# |  |  |  .------- month (1 - 12) OR jan,feb,mar,apr ...
# |  |  |  |  .---- day of week (0 - 6) (Sunday=0 or 7) OR sun,mon,tue,wed,thu,fri,sat
# |  |  |  |  |
# *  *  *  *  * user-name  command to be executed

# After reboot, wait 15 seconds and start cloud services
# -----------------------------------------------------------------------------------
  @reboot sleep 15 && ~/itops-scripts/virtual-machines/cloud/start-cloud.sh

# Virtual machine operating system update (Every Sunday at 23:00 PM)
# -----------------------------------------------------------------------------------
  0  23 *  *  sun   ~/itops-scripts/virtual-machines/cloud/update-os.sh

# Seafile & Mattermost MySQL (MariaDB) databases backup (Every day 01:55 & 01:45 PM)
# -----------------------------------------------------------------------------------
  55 01 *  *  *     ~/itops-scripts/virtual-machines/cloud/backup-seafile-db.sh
  45 01 *  *  *     ~/itops-scripts/virtual-machines/cloud/backup-mattermost-db.sh

# Synchronize IT Operations GitHub repository every 1 hour
# -----------------------------------------------------------------------------------
  35 *  *  *  *     ~/itops-scripts/utils/git-sync.sh ~/itops-scripts/

# Verify every hour, from Monday through Friday, that Seafile cloud services are up
# -----------------------------------------------------------------------------------
  0  *  *  *  1-5   ~/itops-scripts/virtual-machines/cloud/start-cloud.sh

# Seafile storage checking and repair (Every Saturday at 7:00 AM)
# -----------------------------------------------------------------------------------
  0  7  *  *  sat   ~/itops-scripts/virtual-machines/cloud/seafile-storage-cleanup.sh

