# cron/hta-opportunity.crontab
# Example of job definition:
# .---------------- minute (0 - 59)
# |  .------------- hour (0 - 23)
# |  |  .---------- day of month (1 - 31)
# |  |  |  .------- month (1 - 12) OR jan,feb,mar,apr ...
# |  |  |  |  .---- day of week (0 - 6) (Sunday=0 or 7) OR sun,mon,tue,wed,thu,fri,sat
# |  |  |  |  |
# *  *  *  *  * user-name  command to be executed

# After reboot, wait 120 seconds and start all virtual machines
# -----------------------------------------------------------------------------------
@reboot sleep 120 && ~/itops-scripts/physical-machines/vboxall-on.sh

# Server operating system update (Every Tuesday at 21:30 PM)
# -----------------------------------------------------------------------------------
30 21 *  *  tue   ~/itops-scripts/physical-machines/update-os.sh

# VirtualBox full virtual machines backup (Every Sunday at 07:00 AM)
# -----------------------------------------------------------------------------------
0  7  *  *  sun   ~/itops-scripts/physical-machines/vboxbackup.sh

# Synchronize IT Operations GitHub repository every 1 hour
# -----------------------------------------------------------------------------------
30 *  *  *  *     crontab -l > ~/itops-scripts/cron/$(hostname -s).crontab
35 *  *  *  *     ~/itops-scripts/utils/git-sync.sh ~/itops-scripts/

# Verify internet every Every 15 minutes, starting 5 minutes past the hour, Monday through Saturday.
# On Sundays, every 12 hours, starting at 04:00 AM, to allow VMs backup.
# -----------------------------------------------------------------------------------
#5-59/15 * * * 1-6 ~/itops-scripts/physical-machines/net-monitor.sh &>> /data/VirtualBox-VMs/VMLogs/Net-Services.log
#0 4-23/12  * * 0   ~/itops-scripts/physical-machines/net-monitor.sh &>> /data/VirtualBox-VMs/VMLogs/Net-Services.log
# Clear Net Services log every 3 months.
#0  0  * */3 *     echo "" > /data/VirtualBox-VMs/VMLogs/Net-Services.log
# If running, stop the firewall VM to let the main communication host run it.
#57 22 *  *  mon   ~/itops-scripts/physical-machines/internet-off.sh
#30 4  *  *  *     ~/itops-scripts/physical-machines/internet-off.sh

# Export Firewall settings if running on this server (Every day at 03:50 AM)
# -----------------------------------------------------------------------------------
#50 3  *  *  *     ~/itops-scripts/physical-machines/firewall-sync.sh

# Start non-backed-up virtual machines (Every day at 23:30 PM)
# -----------------------------------------------------------------------------------
30 23 *  *  *     ~/itops-scripts/physical-machines/vm-on.sh HTA-ALT

# Restart HTA-RDS remote desktop server (After host scheduled shutdowns)
# -----------------------------------------------------------------------------------
# 0  3  *  *  mon   ~/itops-scripts/physical-machines/vm-restart.sh HTA-RDS
# 0  3  *  *  thu   ~/itops-scripts/physical-machines/vm-restart.sh HTA-RDS

# -----------------------------------------------------------------------------------
0 */4  * * *     ~/itops-scripts/physical-machines/vm-on.sh "HTA-BKP"
#*/30 * * * *     ~/itops-scripts/physical-machines/vm-on.sh "Testbed-01"
#*/35 * * * *     ~/itops-scripts/physical-machines/vm-on.sh "Testbed-02"
#*/40 * * * *     ~/itops-scripts/physical-machines/vm-on.sh "Testbed-03"

